# etc/supervisor-kea-dhcp4.conf
[program:kea-dhcp4]
user=root
command=/bin/start-kea-dhcp4
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
