# etc/supervisor-ash.conf
[program:ash]
user=root
command=/bin/ash -i
autostart=true
autorestart=true
